# config/crontab
42 0 1 * * letsencrypt renew >/dev/null 2>&1
