# config/crontab
42 0 * * * letsencrypt renew >/dev/null 2>&1
