# config/crontab
42 0 1 * * root letsencrypt renew >/dev/null 2>&1
